# roles/podman/it-tools/templates/it-tools.container.j2
[Container]
ContainerName={{ container_name }}
Image={{ container_image }}
NoNewPrivileges=true

[Service]
Restart=always

[Install]
WantedBy=default.target
