# roles/podman/it-tools/files/ittools.container
[Container]
ContainerName=it-tools
Image=docker.io/corentinth/it-tools:latest
PublishPort=8282:80

[Install]
WantedBy=default.target
