# roles/apps/restic/templates/restic-b2.timer.j2
[Unit]
Description=Scheduling restic backup

[Timer]
OnCalendar={{ restic_backup_schedule }}
Persistent=true

[Install]
WantedBy=timers.target

