# roles/podman/quadlet/files/ittools.container
[Container]
ContainerName={{ container_name }}
Image={{ container_image }}

[Service]
Restart=always

[Install]
WantedBy=default.target
